# rescan.service
[Unit]
Description=PCI rescan

[Service]
Type=oneshot
ExecStart=/opt/pci/rescan.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target